# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:vpn_script]
command=/bin/bash /setup.sh
